# docker/etc/s6-overlay/s6-rc.d/plexripper/run
#!/usr/bin/with-contenv sh
exec dotnet /app/PlexRipper.WebAPI.dll